# linux/etc/getty@tty5.service.d/override.conf
[Service]
ExecStart=
Type=simple
ExecStart=-/sbin/agetty --autologin root --noclear %I $TERM
