# linux/etc/getty@tty1.service.d/override.conf
[Service]
ExecStart=
Type=simple
ExecStart=-/sbin/agetty --autologin andy --noclear %I $TERM
